# etc/supervisord.d/neo4j.ini
[program:neo4j]
command=/opt/qnib/bin/start_neo4j.sh
redirect_stderr=true
stdout_logfile=/var/log/supervisor/%(program_name)s.log
autorestart = false
autostart = true
startsecs=2
exitcodes=0
